# celery-kanari.service
# =============================================================================
# Kanari Agent - Systemd Service
# =============================================================================
# Install at: /etc/systemd/system/celery-kanari.service
#
# Commands:
#   sudo systemctl daemon-reload
#   sudo systemctl enable celery-kanari
#   sudo systemctl start celery-kanari
#   sudo journalctl -u celery-kanari -f
# =============================================================================

[Unit]
Description=Kanari Agent - Silent Queue Guardian
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/kanari-agent
ExecStart=/opt/kanari-agent/venv/bin/python kanari_agent.py --config /opt/kanari-agent/config.yaml
Restart=always
RestartSec=10

# Environment
Environment=PYTHONUNBUFFERED=1

# Optional: Environment variables in separate file
# EnvironmentFile=/opt/kanari-agent/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=celery-kanari

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/kanari-agent

[Install]
WantedBy=multi-user.target
